# crontab
#run python script daily at midnight
0 0 * * * python /opt/asx/src/main.py > /proc/1/fd/1 2>/proc/1/fd/2

